# tmp_directories/tmp_dir.service
[Unit]
Description=handles own tmp directories

[Service]
Type=oneshot
ExecStart=/home/cardoso/.tmp_dir
